# SICK_Program/SICK.service
[Unit]
Description=SICK scanners

[Service]
Type=notify
ExecStart=/usr/sbin/SICK_Program

[Install]
WantedBy=multi-user.target